# assets/debian/systemd/stonenetd.service
[Unit]
Description=Stonenet

[Service]
ExecStart=/usr/bin/stonenetd
